# manito/etc/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn.sock
SocketUser=ubuntu

[Install]
WantedBy=sockets.target